# config/feeldscope-webapp.service
[Unit]
Description=FEELDSCOPE Web Application
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/FEELDSCOPE/webapp
ExecStart=/usr/local/bin/npx next start -p 80
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
AmbientCapabilities=CAP_NET_BIND_SERVICE
StandardOutput=journal
StandardError=journal
SyslogIdentifier=feeldscope-webapp

[Install]
WantedBy=multi-user.target
